# com.local.yfinance-movers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.yfinance-movers</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sandarbh/.local/bin/uv</string>
        <string>run</string>
        <string>app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sandarbh/Projects/github/yfinance-movers</string>

    <!-- Starts automatically when you log in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restarts automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logs, so you can check what happened if something goes wrong -->
    <key>StandardOutPath</key>
    <string>/tmp/yfinance-movers.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/yfinance-movers.error.log</string>
</dict>
</plist>
